# coreos/spm-agent-v2.service
[Unit]
Description=SPM Docker Agent
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
Restart=always
RestartSec=30s

ExecStartPre=-/bin/sh -c "/bin/etcdctl get /sematext.com/myapp/logsene/gateway_port || /bin/etcdctl set /sematext.com/m
yapp/logsene/gateway_port 9000"
ExecStartPre=/bin/sh -c "/usr/bin/docker kill spm-agent; true;"
ExecStartPre=/bin/sh -c "/usr/bin/docker rm spm-agent; true;"
ExecStartPre=/bin/sh -c "/usr/bin/docker pull sematext/spm-agent-docker;true"
ExecStart=/bin/sh -c "set -ex; /usr/bin/docker run --name spm-agent \
                      -e SPM_TOKEN=$(etcdctl get /sematext.com/myapp/spm/token) \
                      -e HOSTNAME=$HOSTNAME -v /var/run/docker.sock:/var/run/docker.sock \
                      -e LOGSENE_TOKEN=$(etcdctl get /sematext.com/myapp/logsene/token) \
                      -p $(etcdctl get /sematext.com/myapp/logsene/gateway_port):9000 \
                      sematext/spm-agent-docker"

ExecStop=/usr/bin/docker stop spm-agent

[Install]
WantedBy=multi-user.target

[X-Fleet]
Global=true
